# src/holoiso/etc/systemd/system/sk-expand-home.service
[Unit]
Description=Home Partition Resize
After=network.target

[Service]
#Type=oneshot
ExecStart=/usr/bin/expand_home_partition.expect
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target