# src/chimeraos-addon/lib/udev/rules.d/90-x3d-vcache.rules
# Allow user-space access to AMD X3D V-Cache mode switching
ACTION=="add", SUBSYSTEM=="platform", KERNEL=="AMDI0101:*", RUN+="/bin/chmod 0666 /sys%p/amd_x3d_mode"
