# src/chimeraos-addon/lib/udev/rules.d/99-cec-bluetooth.rules
# from https://github.com/ublue-os/bazzite/pull/1105
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/cec-control onboot"