# src/chimeraos-addon/systemd/user/gamescope-session-plus@.service.d/x3d-vcache.conf
[Service]
ExecStartPre=-/usr/bin/sk-x3d-vcache-set cache
ExecStopPost=-/usr/bin/sk-x3d-vcache-set frequency
