# src/chimeraos-addon/systemd/system/hostname-fix.service
[Unit]
Description=Set hostname to alphanumeric characters
After=network.target
After=sk-setup-kernel-options.service

[Service]
Type=oneshot
ExecStart=/usr/bin/hostname-fix
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target