# src/chimeraos-addon/systemd/system/sk-chos-daemon.service
[Unit]
Description=SkorionOS daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sk-chos-daemon

[Install]
WantedBy=multi-user.target