# src/chimeraos-addon/systemd/system/sk-chos-tool-autoupdate.timer
[Unit]
Description=Trigger sk-chos-tool updates daily and at boot

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
