# src/chimeraos-addon/systemd/user/desktop-monitor.service
[Unit]
Description=Desktop Permissions Monitor

[Service]
Type=simple
ExecStart=/usr/bin/sk-desktop-monitor

[Install]
WantedBy=default.target